# omega-backup-1.0/omega-backup.cron
# Omega Backups for MySQL:
0 0 * * Sun /usr/bin/omega-backup --full > /dev/null 2>&1
0 0 * * Mon-Sat /usr/bin/omega-backup --incr > /dev/null 2>&1
